# root/etc/services.d/jupyterlab/run
#!/usr/bin/zsh
exec s6-setuidgid user /runjupyter.sh

