# root/etc/services.d/jupyterlab/run
#!/usr/bin/with-contenv /bin/zsh
exec s6-setuidgid "$NB_USER" /runjupyter.sh

